# service/ai-wallpaper.timer
[Unit]
Description=Run AI wallpaper generator every day

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
