# service/daily-wall.timer
[Unit]
Description=Run daily wallpaper every day

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
